# crontab.cron
*/5 * * * * cd /home/pi/timelapse && . .venv/bin/activate && ./main.py >> /dev/null
0 1 * * * cd /home/pi/timelapse && .venv/bin/activate && ./fetch.py >> /dev/null
